# data/ayatana-indicators-pre.target
[Unit]
Description=Target representing the lifecycle of the Ayatana Indicators. Each indicator should be bound to it in its individual service file.
